# systemd/airsoft-prop.service
[Unit]
Description=Airsoft Prop
After=multi-user.target sound.target
Wants=network-online.target sound.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=__INSTALL_DIR__
# Give USB audio devices time to enumerate before starting the app.
ExecStartPre=/bin/sleep 3
ExecStart=__VENV_DIR__/bin/python -m src.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Graceful shutdown timeout: give app 15s to clean shutdown via SIGTERM
# before systemd force-kills with SIGKILL
TimeoutStopSec=15

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=GPIOZERO_PIN_FACTORY=rpigpio

[Install]
WantedBy=multi-user.target
